# shadowsocks-server/files/supervisor/obfsproxy_cleanup.ini
[program:obfsproxy_cleanup]
command=/bin/bash /workaround_purgeobfsdata.sh
numprocs=1                    ; number of processes copies to start (def 1)
directory=/tmp                ; directory to cwd to before exec (def no cwd)
umask=022                     ; umask for process (default None)
autostart=true                ; start at supervisord start (default: true)
startsecs=3                   ; number of secs prog must stay running (def. 1)
startretries=5                ; max # of serial start failures (default 3)
stopsignal=KILL
user=root
